# etc/supervisord.d/init.ini
[program:extension-init]
process_name = extension-init
command=/init.sh
startsecs=0
autorestart=false
priority=0
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0