# etc/supervisord.d/postfix.ini
[program:postfix]
;process_name = master
command=/postfix-service.sh
autorestart=unexpected
startretries=10

priority=10
