# etc/supervisord.d/postfix.ini
[program:postfix]
process_name = master
command=/usr/sbin/postfix -c /etc/postfix start
autorestart=false
priority=10
