# etc/supervisord.d/s3-config.ini
[program:s3-config]
process_name = s3-sync
command=/s3-config.sh
startsecs=0
autorestart=false
priority=0
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0